# config/systemd/user/bing-wallpaper-daily.timer
[Unit]
Description=Run wallpaper downloader daily and on boot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target