# config/systemd/user/caffeine-as.service
[Unit]
Description=Starts the caffeine applet
[Service]
ExecStart=/usr/bin/caffeine