# config/systemd/user/clipit-as.service
[Unit]
Description=Starts clipit applet

[Service]
ExecStart=/usr/bin/clipit