# config/systemd/user/disable-caps.timer
[Unit]
Description=keeps caps disabled
Requires=disable-caps.service

[Timer]
OnBootSec=30s
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target