# config/systemd/user/disable-caps-as.service
[Unit]
Description=disables the caps key on start

[Service]
ExecStart=~/disable-caps.sh
