# config/systemd/user/wallpaper-timer.service
[Unit]
Description=<service description>

[Service]
ExecStart=/home/erin/wallpaper-timer.sh 7m
