# config/systemd/user/wallpaper-timer.service
[Unit]
Description=<service description>

[Service]
ExecStart=~/wallpaper-timer.sh 7m
